# launchd/com.quant.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quant.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/zchtech002/machengyu/quant_multifactor/scheduler.py</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zchtech002/machengyu/quant_multifactor</string>

    <key>StandardOutPath</key>
    <string>/Users/zchtech002/machengyu/quant_multifactor/logs/scheduler_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zchtech002/machengyu/quant_multifactor/logs/scheduler_launchd_err.log</string>

    <!-- 登录后自动启动 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 崩溃后自动重启 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 防止 macOS 休眠时杀进程 -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
